# toys.cron
*/5 * * * * cd /home/ubuntu/toys && bash data_update.sh >> /tmp/out.log 2>&1
